# backup/backup.service
[Unit]
Description=Restic backup to Hetzner Storage Box
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
User=admin
ExecStart=/home/admin/docker/backup/backup.sh
Environment=HOME=/home/admin
Nice=19
IOSchedulingClass=idle

[Install]
WantedBy=multi-user.target
