# lando-deploy/lando-deploy-webhook.service
[Unit]
Description=Lando GitHub OIDC Deploy Webhook
After=network-online.target
Wants=network-online.target

[Service]
User=admin
Group=admin
WorkingDirectory=/home/admin/lando-deploy-webhook
ExecStart=/usr/bin/python3 /home/admin/lando-deploy-webhook/server.py
Restart=always
RestartSec=3
Environment=PYTHONUNBUFFERED=1
Environment=LANDO_DEPLOY_CONFIG=/home/admin/lando-deploy-webhook/config.json
EnvironmentFile=-/home/admin/lando-deploy-webhook/secrets.env

[Install]
WantedBy=multi-user.target
